# .crontab
# load vipa workflows samples [once in every two hour]
0 */2 * * * [PROJECT_PATH]/app/console vipa:normalize:journal:workflow